# tools/config-samples/systemd/kapp-backup.service
[Unit]
Description=Schedule of a backup of the k-app database

[Service]
Type=oneshot
WorkingDirectory=/srv/kapp
ExecStart=npm run cli backup
