# tools/config-samples/systemd/kapp-backup.timer
[Unit]
Description=Timer for daily backup of %i

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
